# rpm/after-remove.sh
#!/bin/sh
set -e
chkconfig --del instrument_server
exit 0
