# rpm/after-remove.sh
#!/bin/sh
set -e
chkconfig --del instrumentald
exit 0
