# debian/after-remove.sh
#!/bin/sh
set -e
update-rc.d -f instrumentald remove
exit 0
