# debian/after-remove.sh
#!/bin/sh
set -e
update-rc.d -f instrument_server remove
exit 0
